# easyos/easy-code/rootfs-skeleton/etc/profile.d/asound
#200109 easy buster 2.2.2, bluetooth audio, with bluez-alsa, needs this...
export LIBASOUND_THREAD_SAFE=0
