# easyos/easy-code/rootfs-skeleton/etc/profile.d/ec-update-fix
#20220922 easy 4.4 changed option to run in container as user spot, keep this for
# old containers....
export EC_ENV_ZEUS='false'
