# easyos/easy-code/rootfs-skeleton/etc/profile.d/qemu
#20230717 ref: https://forum.puppylinux.com/viewtopic.php?t=9064
export QEMU_AUDIO_DRV=pa
