# easyos/easy-code/rootfs-skeleton/etc/profile.d/shfm
#20200808 for shfm console file manager...
export SHFM_OPENER=/usr/bin/xdg-open-nox
