# quirky/quirky-code/rootfs-skeleton/etc/profile.d/qt5
export QT5DIR=/usr
#this requires 'qtstyleplugins' to be installed...
export QT_QPA_PLATFORMTHEME=gtk2
# export QT_STYLE_OVERRIDE=gtk2
# export XDG_CURRENT_DESKTOP=GNOME
#menu font is smaller than the gtk2 apps, so do this... nah seems ok...
# export QT_SCALE_FACTOR=1.1
##qt5 compiled in oe, some /usr/lib/cmake .cmake file need this set:
#export OE_QMAKE_PATH_EXTERNAL_HOST_BINS=/usr/bin ...NO. not read in .cmake files.
