# easyos/easy-code/packages-templates/seamonkey/etc/profile.d/seamonkey
export MOZILLA_FIVE_HOME=/usr/lib/seamonkey
