# easyos/easy-code/rootfs-skeleton/usr/local/woofV/rootfs/0rootfs-skeleton/extra/etc/profile.d/xbps-aliases
#run xbps-install via wrapper. note /usr/bin/xbps-install will exec the binary.
alias xbps-install="xbps-install.sh"
#ditto...
alias xbps-remove="xbps-remove.sh"
